# scripts/startup/nlc-cms.service
[Unit]
Description=NLC-CMS Complaint Management System
Documentation=https://github.com/your-org/nlc-cms
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=nlc-cms
Group=nlc-cms
WorkingDirectory=/opt/nlc-cms
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=4005
ExecStart=/usr/bin/node server/server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nlc-cms

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/nlc-cms/logs /opt/nlc-cms/uploads
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target